# docker/browser-stack/supervisord.conf
; Supervisord-Config fuer den Browser-Stack-Container.
; Startet drei Always-On-Services: Xvfb, x11vnc, websockify+noVNC und die Control-API.
; Chromium wird NICHT hier gestartet — die Control-API tut das on demand pro Session.

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 %(ENV_SCREEN_RESOLUTION)s -ac +extension RANDR +extension RENDER -noreset
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

[program:fluxbox]
command=/usr/bin/fluxbox -display :99
environment=DISPLAY=":99"
autostart=true
autorestart=true
priority=15
startsecs=2
stdout_logfile=/var/log/fluxbox.log
stderr_logfile=/var/log/fluxbox.err.log
stdout_logfile_maxbytes=512KB
stderr_logfile_maxbytes=512KB

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -rfbport 5900 -nopw -quiet -noxrecord -noxfixes -noxdamage
autostart=true
autorestart=true
priority=20
startsecs=2
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.err.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 localhost:5900
autostart=true
autorestart=true
priority=30
startsecs=2
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.err.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

[program:control_api]
command=/usr/bin/python3 -m uvicorn control_api:app --host 0.0.0.0 --port 8888 --log-level info
directory=/opt/browser-stack
environment=DISPLAY=":99",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
priority=40
startsecs=3
stdout_logfile=/var/log/control_api.log
stderr_logfile=/var/log/control_api.err.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB
